# etc/supervisord.d/ipynb.conf
[program:ipynb]
command=bash -l -c "/opt/python/bin/ipython notebook --no-browser"
environment=HOME=/home/researcher,USER=researcher,DISPLAY=:99
directory=/home/researcher
autorestart=true
stdout_logfile=/var/log/supervisor/ipynb.log
redirect_stderr=true
